# docker-data/build/supervisor/conf/supervisord.conf
[supervisord]
nodaemon=true
pidfile=./docker-data/log/supervisord.pid
logfile=./docker-data/log/supervisord.log

[program:livestream-manager]
directory=/var/www/livestream/
command=bin/console app:messaging-queue-worker 500
stdout_logfile=./docker-data/log/supervisord-stoutput.log
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
user=livestream
startsecs=2
